# config/upstart-nginx.conf
# nginx.conf

start on runlevel [2345]
stop on runlevel [016]

console owner

respawn

exec /opt/nginx/sbin/nginx -c /opt/nginx/conf/nginx.conf -g "daemon off;"
